# etc/openstack-notifier/openstack-notifier.service
[Unit]
Description=EayunStack Notifier Service
After=network.target

[Service]
ExecStart=/usr/bin/openstack-notifier
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
